# PrimS/PrimS.socket
[Unit]
Description=Socket for Primitier multilayer server


[Socket]
ListenFIFO=/tmp/prims.stdin
Service=PrimS.service